# resources/systemd/timers/locker_config/locker_config.service
[Unit]
Description=Refreshes betterlockscreen backgrounds

[Service]
ExecStart=betterlockscreen -u /usr/share/backgrounds

[Install]
WantedBy=multi-user.target